# saltstack/final/salt/cis-controls/etc/profile.d/timeout.sh
# Ensure default user shell timeout is 900 seconds
readonly TMOUT=900
export TMOUT