# saltstack/base/salt/hostname/etc/systemd/system/google-guest-agent.service.d/override.conf
[Service]
ExecStartPost=/usr/local/bin/cleanup-hosts.sh