# saltstack/base/salt/salt/etc/systemd/system/salt-minion.service
[Unit]
Description=Salt Minion
After=syslog.target network.target

[Service]
Type=simple
Restart=always
ExecStart=/opt/salt/scripts/salt-python-wrapper.sh {{ pillar['salt_path'] }}/bin/salt-minion

[Install]
WantedBy=multi-user.target

[Service]
Environment="SHELL=/bin/bash"
